# app-laptop/thinkfan/files/thinkfan-openrc
#!/sbin/openrc-run

extra_started_commands="reload"

depend() {
	after modules
}

start() {
	ebegin "Starting thinkfan"
	start-stop-daemon --start --exec /usr/sbin/thinkfan -- -q -s5 -c /etc/thinkfan.conf
	eend $?
}

stop() {
	ebegin "Stopping thinkfan"
	start-stop-daemon --stop --exec /usr/sbin/thinkfan
	eend $?
}

reload() {
	PID=$(</var/run/thinkfan.pid)
	ebegin "Sending SIGHUP to thinkfan($PID)"
	kill -HUP $PID
	eend $?
}
